# ulrepo-build.service
[Unit]
Description=ulrepo build services

[Service]
Type=simple
ExecStart=/usr/bin/aursync -c -T
StandardError=journal
